# deploy/monitor/backup/redwall-box-backup.service
[Unit]
Description=Redwall monitör kutusu gece yedeği (pg+sqlite+config)
After=docker.service
Requires=docker.service
[Service]
Type=oneshot
ExecStart=/opt/monitor/backup/box-backup.sh
