# deploy/monitor/deadman/redwall-deadman.timer
[Unit]
Description=Redwall dead-man ping her dakika
[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
[Install]
WantedBy=timers.target
